# haraka.service
[Unit]
Description=Haraka SMTP Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=haraka
Group=haraka
ExecStart=/usr/bin/node haraka.js
WorkingDirectory=/opt/tempmail/haraka
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=haraka
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tempmail/haraka/logs
ReadWritePaths=/opt/tempmail/haraka/queue

# Environment
Environment=NODE_ENV=production
Environment=HARAKA_CONFIG=/opt/tempmail/haraka/config

[Install]
WantedBy=multi-user.target